# init/xsession.target
[Unit]
Description=Xsession user target
BindsTo=graphical-session.target
